# scripts/com.sw.poll-issues.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sw.poll-issues</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shane/pa/main/repos/simple-workflow/scripts/poll-issues.sh</string>
        <string>--once</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>WorkingDirectory</key>
    <string>/Users/shane/pa/main/repos/simple-workflow</string>

    <key>StandardOutPath</key>
    <string>/Users/shane/Library/Logs/simple-workflow/poll.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shane/Library/Logs/simple-workflow/poll.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
